# rtlamr/rootfs/etc/services.d/rtlamr/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Run RTL-AMR
# Add on
# ==============================================================================

bashio::log.info "Starting RTLAMR & RTL Collect"

for line in $(cat '/config/rtlamr-collect.env'); do
    export $line
done

exec /go/bin/rtlamr | /go/bin/rtlamr-collect

Sleep 15