# rtlamr/rootfs/etc/services.d/rtlamr/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Run RTL-AMR
# Add on
# ==============================================================================

# Wait for rtl_tcp to become available
bashio::net.wait_for 1234

# Start RTL-AMR & RTL-COLLECT
bashio::log.info "Starting RTLAMR & RTL Collect"

for line in $(cat '/config/rtlamr-collect.env'); do
    export $line
done

exec /go/bin/rtlamr | /go/bin/rtlamr-collect

sleep 15